# server_setup/grid-db.service
[Unit]
Description=GRID Database (PostgreSQL native)
# PostgreSQL is managed by system postgresql.service
# This unit is a dependency alias so other GRID services can depend on it

After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
