# server_setup/grid-api.service
[Unit]
Description=GRID Intelligence API
After=network.target postgresql.service grid-llamacpp.service
Requires=postgresql.service

[Service]
Type=simple
User=grid
WorkingDirectory=/home/grid/grid_v4/grid_repo/grid
EnvironmentFile=/home/grid/grid_v4/grid_repo/grid/.env
ExecStart=/home/grid/grid_v4/grid_repo/grid/venv/bin/python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/data/grid/logs/api.log
StandardError=append:/data/grid/logs/api.log

[Install]
WantedBy=multi-user.target
