# server_setup/grid-llamacpp.service
[Unit]
Description=GRID llama.cpp Inference Server (Nemotron-Cascade-2 30B-A3B)
After=network.target

[Service]
Type=simple
User=grid
WorkingDirectory=/data/grid_v4/grid_repo
EnvironmentFile=/data/grid_v4/grid_repo/.env
Environment=LD_LIBRARY_PATH=/data/vendor/llama.cpp/build/bin
Environment=LLAMACPP_MODEL=/data/models/nvidia_Nemotron-Cascade-2-30B-A3B-Q4_K_S.gguf
Environment=LLAMACPP_PARALLEL=1
ExecStart=/bin/bash scripts/start_llamacpp.sh
Restart=always
RestartSec=10
StandardOutput=append:/data/grid/logs/llamacpp.log
StandardError=append:/data/grid/logs/llamacpp.log

# GPU access
SupplementaryGroups=video render

[Install]
WantedBy=multi-user.target
